# api/dist/systemd/fdcapi.service
[Unit]
Description=BFPD Branded Food Products REST api service
After=network.target nsqd.service
[Service]
WorkingDirectory=/opt/bfpd
ExecStart=/opt/bfpd/bin/bfpd -c /opt/bfpd/config.yml -p 8001
ExecReload=/bin/kill -HUP $MAINPID
Type=simple
KillMode=process
User=bfpd

[Install]
WantedBy=multi-user.target
